# launchd/com.llm-seo-automation.lg_daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.llm-seo-automation.lg_daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/kcc/Desktop/llm-seo-automation-system && source .venv/bin/activate && export AUTO_RUN_MODE=lg_daily && python llm_seo_system/auto_runner.py</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/llm_seo_automation_lg_daily.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/llm_seo_automation_lg_daily.err.log</string>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
